# skel/mode644/etc/systemd/system/winterman-apt.timer
[Unit]
Description=Update winterman package.

[Timer]
OnUnitActiveSec=5m
Persistent=true

[Install]
WantedBy=timers.target
